# bot.conf
[program:telegrambot]
directory = /home/ubuntu/valiteach
command=python3 app.py

autostart = true
autorestart=true
envirovment=HOME="/home/ubuntu", USER ="ubuntu"

stderr_logfile=/home/ubuntu/valiteach/logfile_err.log
stdout_logfile=/home/ubuntu/valiteach/logfile.log

